# .scripts/udev-rules/10-usb-notify.rules
ACTION=="add", \
	SUBSYSTEM=="usb", \
	ENV{DEVPATH}=="/devices/pci****:**/****:**:**\.*/usb*/*-*:*\.*", \
	ENV{PRODUCT}!="f0d/92/205", \
	RUN+="/usr/bin/sudo -u nerdjp /home/nerdjp/.scripts/usb-notify-connected.sh"

ACTION=="remove", \
	SUBSYSTEM=="usb", \
	ENV{DEVPATH}=="/devices/pci****:**/****:**:**\.*/usb*/*-*:*\.*", \
	ENV{PRODUCT}!="f0d/92/205", \
	RUN+="/usr/bin/sudo -u nerdjp /home/nerdjp/.scripts/usb-notify-disconnected.sh"
